# uca_app.service
[Unit]
Description=UCA App Django Application
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/uca_app
Environment="PATH=/var/www/uca_app/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=uca_project.settings_production"
# Database credentials are hardcoded in settings_production.py
ExecStart=/var/www/uca_app/venv/bin/gunicorn --config /var/www/uca_app/gunicorn.conf.py uca_project.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
